# config/supervisord/supervisord-ui.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
redirect_stderr=true

[program:ui]
command=/usr/local/bin/start_desktop.sh
autostart=true
autorestart=true
redirect_stderr=true
