# pi-esys/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:tftpd]
command=/usr/sbin/in.tftpd  --verbose -L --user tftp --address 0.0.0.0:69 --secure /srv/tftp
